# services/openrc/openrc-service
#!/sbin/openrc-run
description="gimmeasearx - Find a random searx instance"
PORT=':8080'
export PORT
cd 'TEMPLATE_DIR'
pidfile="/run/${RC_SVCNAME}.pid"
command='/usr/local/bin/gimmeasearx'
command_user=nobody
command_background=true

depend() {
	use net
}
